# docker_comm/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:app]
command=/srv/bin/cont/start
